# deployment/packaging/robojacketsstsl.service.in
[Unit]
Description=RoboJackets STSL service unit

[Service]
Type=simple
ExecStart=/home/debian/ros2/eloquent/stsl/lib/stsl/run_node.sh <ROBOT_NUMBER>
User=debian

[Install]
WantedBy=multi-user.target
